# lab/fileserver/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd.log
stdout_logfile=/var/log/supervisor/sshd.log

[program:proftpd]
command=/usr/sbin/proftpd --nodaemon --config /etc/proftpd/proftpd.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/proftpd.log
stdout_logfile=/var/log/supervisor/proftpd.log

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/smbd.log
stdout_logfile=/var/log/supervisor/smbd.log

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nmbd.log
stdout_logfile=/var/log/supervisor/nmbd.log
